# bbu_app/bbu_mgr_app/bbutest&senddata/bbu_send_data_bmc.service
[Unit]
Description=send bbu data to BMC

[Service]
ExecStart=/home/taddeo/fortest_tmp/bbuapp_send_data_BMC
#ExecStop=/home/taddeo/fortest_tmp/bbutest_disarm_bbu
ExecStopPost=/home/taddeo/fortest_tmp/bbutest_disarm_bbu

[Install]
WantedBy=multi-user.target
